# s6-overlay/s6-rc.d/php-install-extensions/up
/etc/s6-overlay/scripts/php-install-extensions