# s6-overlay/s6-rc.d/supercronic-init/up
/etc/s6-overlay/scripts/supercronic-init